# stow/scripts/.config/systemd/user/ags.service
[Unit]
Description=AGS (Aylur's GTK Shell)
Documentation=https://github.com/aylur/ags
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
Environment="PATH=%h/.local/bin/launchers:%h/.local/bin/system:%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/ags run
Restart=on-failure
RestartSec=3
KillMode=mixed
TimeoutStopSec=5

[Install]
WantedBy=graphical-session.target
