# stow/scripts/.config/systemd/user/ags.service
[Unit]
Description=AGS (Aylur's GTK Shell)
Documentation=https://github.com/aylur/ags
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
# Wait for Wayland display to be available
ExecStartPre=/bin/sh -c 'until [ -n "$WAYLAND_DISPLAY" ]; do sleep 0.1; done'
ExecStart=/usr/bin/ags run
Restart=on-failure
RestartSec=3
# Limit restart attempts
StartLimitBurst=5
KillMode=mixed
TimeoutStopSec=5

[Install]
WantedBy=graphical-session.target
