# stow/scripts/.config/systemd/user/ags.service
[Unit]
Description=AGS (Aylur's GTK Shell)
Documentation=https://github.com/aylur/ags

[Service]
Type=simple
ExecStart=/usr/bin/ags run
Restart=on-failure
RestartSec=3
KillMode=mixed
TimeoutStopSec=5
