# stow/scripts/.config/systemd/user/bluetooth.service
[Unit]
Description=Reconnect Bluetooth Devices
After=bluetooth.service
Wants=bluetooth.service
# Wait for AGS to be fully started and stable
After=ags.service
Wants=ags.service

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 1
ExecStart=/home/john/.local/bin/system-bluetooth
StandardOutput=journal
StandardError=journal
# Don't fail if Bluetooth is not available
SuccessExitStatus=0 1

[Install]
WantedBy=default.target

