# stow/bluetooth/.config/systemd/user/bluetooth-auto-switch.service
[Unit]
Description=Auto-switch audio to Bluetooth devices
After=pipewire.service wireplumber.service

[Service]
Type=simple
ExecStart=%h/.local/bin/system-bluetooth-audio-switch
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
